# .crontab
0 2 * * * /generate_bill_status_codes.sh
